# dump1090.conf
start on runlevel [2345]
#start on (local-filesystems and net-device-up IFACE=eth0 and started autosshseanav)
#start on (local-filesystems)
stop on shutdown

respawn
respawn limit 15 5
# Max open files are @ 1024 by default. Bit few.
#limit nofile 32768 32768
#give console time to get going - so we set the baud rate after it doesto 4800
pre-start script
    /bin/echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting adsb" >> /home/ubuntu/dump1090-fa/startadsb.log
    /bin/sleep 5
end script
script
    . /etc/environment
    export FLEETRANGE_ID
    /bin/echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Started adsb" >> /home/ubuntu/dump1090-fa/startadsb.log
    chdir /home/ubuntu/dump1090-fa
    exec /home/ubuntu/dump1090-fa/dump1090 --lat 51.7200648 --lon -2.6844428 --json-location-accuracy 2 --quiet --gain -10 --net --write-json /var/www/html/adsb/data >> /home/ubuntu/dump1090-fa/startadsb.log
end script
